# lib/desktop/kde/plasma-kwin_wayland.service.in
; SPDX-FileCopyrightText: 2023 Roman Gilg <subdiff@gmail.com>
;
; SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=KDE Window Manager
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kwin_wayland_wrapper --xwayland
BusName=org.kde.KWinWrapper
Slice=session.slice
